# packaging/ubuntu/generic/harmony-smp.service
[Unit]
Description=Harmony eDelivery Service Metadata Publisher
After=network.target

[Service]
Type=forking

Environment=JAVA_HOME=/usr/lib/jvm/java-8-openjdk-amd64
Environment=CATALINA_PID=/opt/harmony-smp/temp/harmony-ap.pid
Environment=CATALINA_HOME=/opt/harmony-smp
Environment=CATALINA_BASE=/opt/harmony-smp
Environment=CATALINA_OUT=/var/log/harmony-smp/catalina.out
Environment='CATALINA_OPTS=-Xms256M -Xmx512M -server -XX:+UseParallelGC'
Environment='JAVA_OPTS=-Djava.awt.headless=true -Djava.security.egd=file:/dev/./urandom'
ExecStart=/opt/harmony-smp/bin/startup.sh
ExecStop=/opt/harmony-smp/bin/shutdown.sh

User=harmony-smp
Group=harmony-smp
UMask=0007
RestartSec=10
Restart=always

[Install]
WantedBy=multi-user.target
